# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=service nginx restart

[program:gunicorn]
command=gunicorn __init__:app